# dashboard-ui.service
[Unit]
Description=Dashboard UI Next.js Application
After=network.target

[Service]
Environment=NODE_ENV=production
# This tells the Next.js app that API calls starting with /api/ should go to the same host.
# Nginx will then route cherry-ai.me/api/... to the backend on localhost:8000.
Environment=NEXT_PUBLIC_API_URL=/api 

WorkingDirectory=/opt/cherry-ai/dashboard

# Install dependencies and build if not already done (consider moving build to a separate deploy script step)
ExecStartPre=/usr/bin/npm install --omit=dev # Use --omit=dev or ci for production dependencies
ExecStartPre=/usr/bin/npm run build

ExecStart=/usr/bin/npm run start # This runs `next start` which typically uses port 3000

Restart=always
User=ubuntu # Or a dedicated non-root user if created
Group=ubuntu # Or a dedicated group

# Output logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dashboard-ui

[Install]
WantedBy=multi-user.target 