# mcp_server/scripts/gcp-resources-mcp-server.service
[Unit]
Description=GCP Resources Model Context Protocol (MCP) Server
After=network.target mcp-server.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=vscode
WorkingDirectory=/workspaces/orchestra-main
ExecStart=/usr/bin/python3 gcp_resources_mcp_server.py --host 127.0.0.1 --port 8085
Environment=PYTHONUNBUFFERED=1
Environment=GCP_SNAPSHOT_DIR=/workspaces/orchestra-main/.gcp-snapshots

[Install]
WantedBy=multi-user.target
