# deployment/orchestra-api.service
[Unit]
Description=Orchestra AI Real Agents API
After=network.target redis.service

[Service]
Type=simple
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
